# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true


[program:gunicorn]
command=/usr/local/bin/gunicorn -w 2 --log-level=debug --bind 0.0.0.0:8077 djangoserver.wsgi:application
autostart=true
autorestart=true